# read_sensors.service
[Unit]
Description=Read sensor data and forward payloads over LoRa
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dynx
Group=dynx
WorkingDirectory=/home/dynx/Rocket
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/dynx/Rocket/external/LoRa-RaspberryPi:$PYTHONPATH"
ExecStart=/usr/bin/env python3 -u /home/dynx/Rocket/read_sensors.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
