# read_sensors.service
[Unit]
Description=Read sensor data and forward payloads over LoRa
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mario
Group=mario
WorkingDirectory=/home/mario/rocket_system
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/mario/rocket_system/external/LoRa-RaspberryPi:$PYTHONPATH"
ExecStart=/usr/bin/env python3 -u read_sensors.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
